# builder/files/etc/systemd/system/ttyd.service
[Unit]
Description=ttyd Web Terminal

[Service]
ExecStart=/usr/local/bin/ttyd login

[Install]
WantedBy=multi-user.target
